# server/deploy/supervisord.ini
[program:julolo]
user=root
command=gunicorn -c deploy/prd_julolo_api.py run:app
autostart=true
autorestart=true
stopasgroup=true
directory=/var/www/julolo/server

[program:julolo_master]
user=root
command=gunicorn -c deploy/prd_julolo_admin.py run_admin:app
autostart=true
autorestart=true
stopasgroup=true
directory=/var/www/julolo/server
