# server/deploy/supervisord.ini
[program:hys]
user=root
command=gunicorn -c deploy/prd_hys_api.py run:app
autostart=true
autorestart=true
stopasgroup=true
directory=/var/www/hys/server

[program:hys_master]
user=root
command=gunicorn -c deploy/prd_hys_admin.py run_admin:app
autostart=true
autorestart=true
stopasgroup=true
directory=/var/www/hys/server
